# setup4packetx/services/open5gs-udmd.service
[Unit]
After=open5gs-ausfd.service

[Service]
ExecStart=/usr/local/open5gs/bin/open5gs-udmd

[Install]
WantedBy=default.target
