# setup4packetx/services/open5gs-pcfd.service
[Unit]
After=open5gs-udmd.service

[Service]
ExecStart=/usr/local/open5gs/bin/open5gs-pcfd

[Install]
WantedBy=default.target
