# setup4packetx/services/open5gs-nat.service
[Unit]
After=network-online.target

[Service]
ExecStart=/etc/iptables.sh

[Install]
WantedBy=default.target
