# setup4packetx/services/open5gs-pcrfd.service
[Unit]
After=open5gs-hssd.service

[Service]
ExecStart=/usr/local/open5gs/bin/open5gs-pcrfd

[Install]
WantedBy=default.target
