# setup4packetx/services/open5gs-amfd.service
[Unit]
After=open5gs-scpd.service

[Service]
ExecStart=/usr/local/open5gs/bin/open5gs-amfd

[Install]
WantedBy=default.target
