# setup4packetx/services/open5gs-nssfd.service
[Unit]
After=open5gs-pcfd.service

[Service]
ExecStart=/usr/local/open5gs/bin/open5gs-nssfd

[Install]
WantedBy=default.target
