# setup4packetx/services/open5gs-ausfd.service
[Unit]
After=open5gs-upfd.service

[Service]
ExecStart=/usr/local/open5gs/bin/open5gs-ausfd

[Install]
WantedBy=default.target
