# setup4packetx/services/open5gs-smfd.service
[Unit]
After=open5gs-amfd.service

[Service]
ExecStart=/usr/local/open5gs/bin/open5gs-smfd

[Install]
WantedBy=default.target
