# setup4packetx/services/open5gs-bsfd.service
[Unit]
After=open5gs-nssfd.service

[Service]
ExecStart=/usr/local/open5gs/bin/open5gs-bsfd

[Install]
WantedBy=default.target
