# setup4packetx/services/open5gs-sgwud.service
[Unit]
After=open5gs-sgwcd.service

[Service]
ExecStart=/usr/local/open5gs/bin/open5gs-sgwud

[Install]
WantedBy=default.target
