# setup4packetx/services/open5gs-nrfd.service
[Unit]
After=mongod.service

[Service]
ExecStart=/usr/local/open5gs/bin/open5gs-nrfd

[Install]
WantedBy=default.target
