# setup4packetx/services/open5gs-scpd.service
[Unit]
After=open5gs-nrfd.service

[Service]
ExecStart=/usr/local/open5gs/bin/open5gs-scpd

[Install]
WantedBy=default.target
